# packaging/debian/prerm
#!/bin/bash
set -e

# Pre-removal script for local-rag
# This script runs before package removal

# Stop and disable the service
if systemctl is-active --quiet local-rag.service; then
    systemctl stop local-rag.service
fi

if systemctl is-enabled --quiet local-rag.service; then
    systemctl disable local-rag.service
fi

echo "Local RAG service stopped and disabled."

exit 0